# preview-worker/worker.conf
[supervisord]
nodaemon=true

[program:preview-worker]
process_name=%(program_name)s_%(process_num)02d
command=python $(WORKER_SCRIPT)s
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true