# lumen-worker/lumen-worker.crontab
*       *       *       *       *       php /var/www/application/artisan queue:work --timeout=0