# lumen-worker/lumen-scheduler.crontab
*       *       *       *       *       php /var/www/application/artisan schedule:run
